# .circleci/twihaialert-app-task.timer
[Unit]
Description=twihaialert task execute every 30 minutes

[Timer]
OnCalendar=*-*-* *:0/30:00
 
[Install]
WantedBy=timers.target